# deploy/dawnop-backend.service
# systemd 单元：用 uvicorn 守护 FastAPI 后端，监听 127.0.0.1:8000（仅本机，由 Nginx 反代对外）。
#
# 安装：
#   1) 把后端代码放到 WorkingDirectory（如 /opt/dawnop/backend），在其中建好 .venv 并 pip install -r requirements.txt
#   2) 准备好 backend/.env（七牛密钥、JWT secret 等，勿提交），首次需跑 python scripts/seed_admin.py
#   3) sudo cp deploy/dawnop-backend.service /etc/systemd/system/dawnop-backend.service
#      （按实际路径/用户改下面的 User / WorkingDirectory / ExecStart）
#   4) sudo systemctl daemon-reload
#      sudo systemctl enable --now dawnop-backend
#      sudo systemctl status dawnop-backend     # 查看运行状态
#      journalctl -u dawnop-backend -f          # 跟踪日志

[Unit]
Description=dawnop-site FastAPI backend (uvicorn)
After=network.target

[Service]
Type=simple
User=dawn
Group=dawn
WorkingDirectory=/opt/dawnop/backend
# .env 由 EnvironmentFile 加载（也可让 app.config 自行读取 backend/.env，二选一即可）
EnvironmentFile=/opt/dawnop/backend/.env
ExecStart=/opt/dawnop/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
